# sunshine-gamescope.path
[Unit]
Description=Watch for Sunshine trigger files

[Path]
PathModified=/home/deck/.trigger-gamescope-set
PathModified=/home/deck/.trigger-gamescope-unset

[Install]
WantedBy=default.target
