# sunshine-gamescope.service
[Unit]
Description=Run Sunshine Gamescope Script

[Service]
Type=oneshot
ExecStart=/home/deck/manage-gamescope.sh
